# generic/lib/systemd/system/nfs.target
[Unit]
Description=NFS

[Install]
WantedBy=multi-user.target
