# generic/lib/systemd/system/pure-uploadscript.service
[Unit]
Description=pure-ftpd upload script helper

[Service]
ExecStart=/usr/sbin/pure-uploadscript @ARGS@ -r @SCRIPT@
