# generic/lib/systemd/system/rc-once.target
[Unit]
Description=System setup
DefaultDependencies=no
After=rc-once.service
Requires=rc-once.service

[Install]
Alias=default.target
