# generic/lib/systemd/system/telnetd.socket
[Unit]
Conflicts=telnetd.service
ConditionKernelCommandLine=ptxdist-devel

[Socket]
ListenStream=23
Accept=yes

[Install]
WantedBy=sockets.target
