# generic/lib/systemd/system/urshd.service
[Unit]
Description=urshd Server
After=syslog.target

[Service]
ExecStart=/usr/bin/urshd -D

[Install]
WantedBy=multi-user.target
