# generic/lib/systemd/system/utelnetd.service
[Unit]
Description=Telnet Server
After=syslog.target

[Service]
ExecStart=/sbin/utelnetd

[Install]
WantedBy=multi-user.target
