# projectroot/usr/lib/systemd/system/logrotate.service
[Unit]
Description=logrotate

[Service]
ExecStart=/usr/sbin/logrotate /etc/logrotate.conf
