# projects/innokom/etc/inittab
# 
# /etc/inittab
#
# Copyright (C) 2003 Auerswald GmbH & Co. KG
# Kai-Uwe Bloem <linux-development@auerswald.de>
#

::sysinit:/etc/init.d/banner
::sysinit:/etc/init.d/mount-proc
::sysinit:/etc/init.d/unlock-flash
::sysinit:/etc/init.d/mount-local
::sysinit:/etc/init.d/root-remount
::sysinit:/opt/auerswald/scripts/data-setup
::sysinit:/etc/init.d/sys-config
::sysinit:/etc/init.d/syslog
::sysinit:/etc/init.d/network-eth
::sysinit:/etc/init.d/network-usb
::sysinit:/etc/init.d/network-ppp
::sysinit:/etc/init.d/mount-nfs
::sysinit:/etc/init.d/read-hwclock
::sysinit:/etc/init.d/net-servers

# Start userland network services 
#::respawn:/sbin/utelnetd
#::respawn:/sbin/proftpd -n

# put a getty on a serial line (for a terminal)
#::respawn:/sbin/getty -L ttyS0 9600 vt100

# Start an "askfirst" shell on the console
::askfirst:-/bin/sh

# Stuff to do before rebooting
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r

# Auerswald
::sysinit:/opt/auerswald/lk.sh
::respawn:opt/auerswald/stats >/dev/null 2>/dev/null </dev/null
::respawn:opt/auerswald/scripts/run-webserver >/dev/null 2>/dev/null </dev/null
::respawn:opt/auerswald/scripts/run-ets >/dev/null 2>/dev/null </dev/null
::respawn:opt/auerswald/scripts/run-syslogd >/dev/null 2>/dev/null </dev/null
