# etc/generic/inittab
# 
# /etc/inittab
#
# $Id: inittab,v 1.6 2004/08/23 08:10:43 rsc Exp $
#

console::sysinit:/etc/init.d/rcS
console::respawn:/sbin/getty -L @SPEED@ @CONSOLE@ vt100

# Stuff to do before rebooting
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r

