# projectroot/lib/systemd/system/collectd.service
[Unit]
Description=statistics collection daemon

[Service]
ExecStart=/usr/sbin/collectd -f
Restart=always

[Install]
WantedBy=multi-user.target
