# projectroot/lib/systemd/system/ntpd.service
[Unit]
Description=ntpd

[Service]
ExecStart=/usr/sbin/ntpd -n -c /etc/ntp-server.conf

[Install]
WantedBy=multi-user.target
