# projectroot/lib/systemd/system/openct.service
[Unit]
Description=OpenCT

[Service]
Type=forking
ExecStart=/usr/sbin/openct-control init
ExecStop=/usr/sbin/openct-control shutdown
