# projectroot/lib/systemd/system/pure-ftpd@.service
[Unit]
Description=pure-ftpd
After=syslog.target
@SCRIPT_DEPS@
[Service]
ExecStart=/usr/sbin/pure-ftpd @ARGS@
StandardInput=socket
