# projectroot/lib/systemd/system/rngd.service
[Unit]
Description=Hardware RNG Daemon
DefaultDependencies=no
Before=local-fs.target

[Service]
ExecStart=/usr/sbin/rngd -f
