# projectroot/lib/systemd/system/var-lib-nfs.mount
[Unit]
Before=local-fs.target

[Mount]
What=tmpfs
Where=/var/lib/nfs
Type=tmpfs
Options=nosuid,nodev,noexec,mode=0755,size=1M
